# networking/aruba/backup-aruba.timer
[Unit]
Description=aruba Backups

[Timer]
OnBootSec=0
OnUnitActiveSec=1days

[Install]
WantedBy=timers.target
