# networking/mikrotik/backup-mikrotik.timer
[Unit]
Description=Mikrotik Backups

[Timer]
OnBootSec=0
OnUnitActiveSec=1days

[Install]
WantedBy=timers.target